# supernode/bsn-supernode-environment.path
[Unit]
Description=Restart supernode if the environment changes

[Path]
PathChanged=/etc/environment

[Install]
WantedBy=multi-user.target
